# systemd/premonition-probe-premarket.timer
[Unit]
Description=Fire the premonition pre-market latency probe Mon-Fri at 08:00 ET

[Timer]
OnCalendar=Mon-Fri 08:00 America/New_York
Persistent=true
AccuracySec=30s

[Install]
WantedBy=timers.target
