# systemd/premonition-probe-premarket.timer
[Unit]
Description=Fire the premonition pre-market latency probe Mon-Fri at 07:55 ET

[Timer]
OnCalendar=Mon-Fri 07:55 America/New_York
Persistent=true
AccuracySec=30s

[Install]
WantedBy=timers.target
