# systemd/premonition-probe-premarket.service
[Unit]
Description=premonition Probe A (freshness) — the original pre-market run, 08:00-09:30 ET, testing whether preMarketVolume ever populates during a real PRE state.

[Service]
Type=oneshot
WorkingDirectory=%h/code/premonition
EnvironmentFile=-/etc/premonition/env
ExecStart=%h/code/premonition/.venv/bin/python3 %h/code/premonition/bin/premonition-probe --probe a --start 08:00 --end 09:30
TimeoutStartSec=infinity
